# platform/meta-birdycamcam/recipes-core/birdycamcam/files/birdycamcam-altroot.service
[Unit]
Description=Configure alternate root partition for Birdy Cam Cam

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@LIBEXECDIR@/birdycamcam/configure-altroot-alias.sh

[Install]
WantedBy=multi-user.target
