# platform/meta-birdycamcam/recipes-multimedia/motion/motion/motion.service
[Unit]
Description=Motion detection video capture daemon
After=network.target

[Service]
ExecStart=motion
Restart=on-failure

[Install]
WantedBy=multi-user.target
