# platform/meta-birdycamcam/recipes-core/autoreboot/autoreboot/autoreboot.service
[Unit]
Description=Reboot the board every day

[Service]
Type=oneshot
ExecStart=reboot
